# roles/alertmanager/templates/alertmanager.service.j2
# This file is managed by Ansible. Manual changes will be overwritten.

[Unit]
Description=Alertmanager for Prometheus
Documentation=https://prometheus.io/docs/alerting/latest/alertmanager/
After=network-online.target
Wants=network-online.target

[Service]
User={{ alertmanager_user }}
Group={{ alertmanager_group }}
Type=simple
Restart=on-failure
RestartSec=5s
ExecStart={{ alertmanager_install_dir }}/alertmanager \
  --config.file={{ alertmanager_config_dir }}/alertmanager.yml \
  --storage.path={{ alertmanager_data_dir }} \
  --web.listen-address=:{{ alertmanager_port }} \
  --log.level=info
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
